# config./systemd/user/openclaw-gateway.service.d/env.conf
[Service]                                                                          
EnvironmentFile=%h/.config/environment.d/openclaw.conf
